# infra/scripts/payment-server.service
[Unit]
Description=Payment Portal Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/payment-portal/server
ExecStart=/usr/local/bin/bundle exec rackup config.ru -p 9292 -o 0.0.0.0
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target